# dist/background-fission.service
[Unit]
Description=background-fission background manager

[Service]
Type=simple
ExecStart=%h/.local/bin/background-fission

[Install]
WantedBy=default.target